# porcupine2/rootfs/etc/s6-overlay/s6-rc.d/porcupine1/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Porcupine 2 service
# ==============================================================================
flags=()

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

exec python3 /usr/src/server.py \
    --uri 'tcp://0.0.0.0:10400' \
    --sensitivity "$(bashio::config 'sensitivity')" \
    --access_key "$(bashio::config 'access_key')" ${flags[@]}
